# etc/crontab
0 3 * * * cd /home/vas3k/vas3k.club && /usr/local/bin/pipenv run python manage.py renew_subscriptions >/home/vas3k/crons.log 2>&1
0 8 * * 2-5 cd /home/vas3k/vas3k.club && /usr/local/bin/pipenv run python manage.py send_daily_digest --production true >/home/vas3k/crons.log 2>&1
0 10 * * 1 cd /home/vas3k/vas3k.club && /usr/local/bin/pipenv run python manage.py send_weekly_digest --production true >/home/vas3k/crons.log 2>&1
0 4 * * * sudo certbot renew
